# .github/runners/airstack-runner.service
[Unit]
Description=AirStack GitHub Actions Runner (ephemeral)
Documentation=https://docs.github.com/en/actions/hosting-your-own-runners
After=network-online.target docker.service
Wants=network-online.target

[Service]
User=runner
Group=runner
WorkingDirectory=/opt/actions-runner

# Place runtime configuration here (REPO_URL, RUNNER_LABELS, etc.).
# Each line: KEY=value  (no quotes needed, no export keyword)
EnvironmentFile=/etc/github-runner-env

ExecStart=/opt/actions-runner/register-runner.sh

# Restart unconditionally so the loop survives transient API errors or reboots.
Restart=always
RestartSec=5

# Give the runner enough time to finish a long job before systemd kills it.
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
